# deploy/supervisor.conf
[program:tuistudio]
command=/var/www/tuistudio/venv/bin/gunicorn -c deploy/gunicorn.conf.py tuistudio:app
directory=/var/www/tuistudio
user=root
autostart=true
autorestart=true