# deploy/supervisor.conf
[program:tm2]
command=/var/www/tuistudio/venv/bin/gunicorn -c deploy/gunicorn.conf tuistudio:app
directory=/var/www/tuistudio
user=root
autostart=true
autorestart=true